# config/launchd/com.example.throughline.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Template: replace {{...}} placeholders, then
     launchctl bootstrap gui/$(id -u) <path-to-this-plist>
     launchctl enable   gui/$(id -u)/com.example.throughline.sync
     launchctl kickstart -k gui/$(id -u)/com.example.throughline.sync

     Purpose: periodically pull conversations exported by the upstream
     OpenWebUI host into the local raw-conversation tree, so the refine
     daemon can pick them up. Default interval is 5 minutes.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.throughline.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{PYTHON}}</string>
        <string>{{THROUGHLINE_HOME}}/scripts/sync_raw_conversations.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{THROUGHLINE_HOME}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>

        <!-- Upstream exporter host (rsync target). Example: "user@host:/srv/openwebui_staging/" -->
        <key>THROUGHLINE_SYNC_SOURCE</key>
        <string>{{SYNC_SOURCE}}</string>

        <!-- Local raw-conversation destination -->
        <key>THROUGHLINE_RAW_ROOT</key>
        <string>{{THROUGHLINE_RAW_ROOT}}</string>
    </dict>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/{{USER}}/Library/Logs/throughline/sync.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/{{USER}}/Library/Logs/throughline/sync.err.log</string>
</dict>
</plist>
